# launchd/com.vadim.jobapply-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Daily job-apply digest. Each morning at 08:00 it runs a live last-24h scan,
  tailors résumé + cover-letter PDFs for newly-ready roles, and emails you the
  apply digest with the PDFs attached. It does NOT submit applications.

  Requires SMTP_* + DIGEST_TO in the repo .env (Gmail app password).

  Install (edit the paths below to your checkout first if different):
    cp launchd/com.vadim.jobapply-daily.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.vadim.jobapply-daily.plist
  Uninstall:
    launchctl bootout gui/$(id -u)/com.vadim.jobapply-daily
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vadim.jobapply-daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/vadimkoenen/Documents/Application Agent/app_main.py</string>
    <string>--auto</string>
    <string>--hours</string>
    <string>24</string>
    <string>--live</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/vadimkoenen/Documents/Application Agent</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/jobapply-daily.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/jobapply-daily.log</string>
</dict>
</plist>
